# systemd/user/gtex62-github-traffic.timer
[Unit]
Description=Periodic refresh of gtex62-core GitHub traffic cache

[Timer]
OnBootSec=10m
OnUnitActiveSec=12h
Persistent=true

[Install]
WantedBy=timers.target
